# conf/invoiceninja.cron
0 8 * * * __APP__ /usr/bin/php__PHP_VERSION__ __INSTALL_DIR__/artisan ninja:send-invoices -q
0 8 * * * __APP__ /usr/bin/php__PHP_VERSION__ __INSTALL_DIR__/artisan ninja:send-reminders -q

